# unifi_radius_retransmit/deploy/systemd/radius.timer
[Unit]
Description=Start UniFi Radius Remediation at boot

[Timer]
OnBootSec=30s
Unit=radius.service

[Install]
WantedBy=timers.target
