# unifi_radius_retransmit/deploy/systemd/radius.service
[Unit]
Description=UniFi Radius Remediation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/radius
Environment=PYTHONPATH=/opt/radius
ExecStart=/opt/radius/.venv/bin/python -m app.scheduler
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
